# vendor/k8s.io/contrib/init/upstart/kube-apiserver.conf
description "Kubernetes API Server"

start on (started networking and started etcd)
stop on runlevel [016]

respawn
setuid kube

script
	if [ -f /etc/kubernetes/config ]; then
	   . /etc/kubernetes/config
	fi
	if [ -f /etc/kubernetes/apiserver ]; then
	   . /etc/kubernetes/apiserver
	fi
	exec /usr/bin/kube-apiserver \
	    $KUBE_LOGTOSTDERR \
	    $KUBE_LOG_LEVEL \
	    $KUBE_ETCD_SERVERS \
	    $KUBE_API_ADDRESS \
	    $KUBE_API_PORT \
	    $KUBELET_PORT \
	    $KUBE_ALLOW_PRIV \
	    $KUBE_SERVICE_ADDRESSES \
	    $KUBE_ADMISSION_CONTROL \
	    $KUBE_API_ARGS
end script
