# vendor/k8s.io/contrib/init/upstart/kube-scheduler.conf
description "Kubernetes Scheduler Plugin"

start on (started kube-apiserver)
stop on runlevel [016]

respawn
setuid kube

script
	if [ -f /etc/kubernetes/config ]; then
	   . /etc/kubernetes/config
	fi
	if [ -f /etc/kubernetes/scheduler ]; then
	   . /etc/kubernetes/scheduler
	fi
	exec /usr/bin/kube-scheduler \
	    $KUBE_LOGTOSTDERR \
	    $KUBE_LOG_LEVEL \
	    $KUBE_MASTER \
	    $KUBE_SCHEDULER_ARGS
end script
