# vendor/k8s.io/contrib/init/upstart/kubelet.conf
description "Kubernetes Kubelet Server"

start on (started docker)
stop on runlevel [016]

respawn

script
	if [ -f /etc/kubernetes/config ]; then
	   . /etc/kubernetes/config
	fi
	if [ -f /etc/kubernetes/kubelet ]; then
	   . /etc/kubernetes/kubelet
	fi
	exec /usr/bin/kubelet \
	    $KUBE_LOGTOSTDERR \
	    $KUBE_LOG_LEVEL \
	    $KUBELET_API_SERVER \
	    $KUBELET_ADDRESS \
	    $KUBELET_PORT \
	    $KUBELET_HOSTNAME \
	    $KUBE_ALLOW_PRIV \
	    $KUBELET_ARGS
end script
