# vendor/k8s.io/contrib/ansible/roles/kubernetes-addons/templates/kube-addons.upstart.j2
description "Kubernetes Addon Object Manager"

start on started kube-apiserver
stop on runlevel [!2345]

env TOKEN_DIR={{ kube_token_dir }}
env KUBECTL_BIN={{ bin_dir }}/kubectl
env KUBERNETES_MASTER_NAME={{ groups['masters'][0] }}
env TRUSTY_MASTER=true

exec {{ kube_script_dir }}/kube-addons.sh
