# vendor/k8s.io/contrib/ansible/roles/contiv/templates/netmaster.service
[Unit]
Description=Netmaster
After=auditd.service systemd-user-sessions.service etcd.service

[Service]
EnvironmentFile=/etc/default/netmaster
ExecStart={{ bin_dir }}/netmaster $NETMASTER_ARGS
KillMode=control-group
