# vendor/k8s.io/contrib/ansible/roles/master/files/delay-master-services.target
[Unit]
Description=Delay Kubernetes services until network is actually online
Wants=network-online.target
After=network-online.target
Before=flanneld.service
Before=kube-apiserver.service

[Install]
WantedBy=flanneld.service
WantedBy=kube-apiserver.service
